# bindfs/etc/supervisor/conf.d/my-bindfs-mounts.conf
[program:my-bindfs-mounts]
command=/usr/local/bin/my-bindfs-mounts --umount-all && /usr/local/bin/my-bindfs-mounts --load-from-file
priority=111

