# network/etc/supervisor/conf.d/my-docker-network.conf
[program:my-docker-network]
command=/usr/local/bin/my-docker-network --start --emit-only -D
startsecs=0
