# network/etc/supervisor/conf.d/my-docker-network.conf
[program:my-docker-network]
command=/usr/local/bin/my-my-docker-network --start -D

