# supervisor/etc/supervisor/conf.d/my-bindfs-mounts.conf
[program:my-bindfs-mounts]
command=/usr/bin/my-bindfs-mounts --load-from-file
priority=111
startsecs=0

